# deploy/systemd/operation-backup.service
[Unit]
Description=Backup Operation bot SQLite database

[Service]
Type=oneshot
WorkingDirectory=@OPERATION_BOT_ROOT@
Environment="MODE=live" "OPERATION_ENV_FILE=@OPERATION_ENV_FILE_LIVE@" "ENV_ROOT=@OPERATION_ENV_ROOT@" "RUN_ROOT=@OPERATION_RUN_ROOT@" "DATA_ROOT=@OPERATION_DATA_ROOT@" "LOG_ROOT=@OPERATION_LOG_ROOT@" "BACKUP_ROOT=@OPERATION_BACKUP_ROOT@"
ExecStart=/usr/bin/env bash @OPERATION_BOT_ROOT@/scripts/backup_sqlite.sh
